# config/.config/autostart/xmousepasteblock.desktop
[Desktop Entry]
Exec=/usr/bin/xmousepasteblock
Icon=dialog-scripts
Name=xmousepasteblock
Path=
Type=Application
X-KDE-AutostartScript=true
